# touchpi/rootfs/home/pi/.local/share/systemd/user/startx.service
[Unit]
Description=startx as user
Wants=network.target

[Service]
Nice=10
Type=simple
Restart=always
## this will automatically start x with default spectrwm session which will start the rest from spectrwm.conf
ExecStart=/usr/bin/startx
SyslogIdentifier=%i

[Install]
WantedBy=default.target
